# 02-Admin/etc/systemd/system/my1.service.d/override.conf
[Service]
LimitCore=infinity
MYSQLD_OPTS="--general_log=0" 
NO_INIT=true
